# etc/supervisor/conf.d/mariadb.conf
[program:mariadb]
command=/usr/sbin/mariadb
user=root
